# deploy/ai-textbook.service
[Unit]
Description=AI-Native Textbook Backend
After=network.target

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/app
EnvironmentFile=/app/.env
ExecStart=/app/deploy/start_production.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-textbook

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/app/logs
ReadWritePaths=/tmp

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target